# supervisord.conf
[supervisord]

[program:websockets=test]
directory = /home/ubuntu/oit-chat-search/
command = python server.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 2
autorestart = true
